# launchopen_loira
 #!/bin/sh
### BEGIN INIT INFO
# Provides:          loira
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: Start/stop loira server
### END INIT INFO

WORK_DIR="/home/das"
NAME="loira"
JAR="loira-0.1-standalone.jar"
USER="das"
DAEMON="/usr/bin/java"
DAEMON_ARGS=" -jar $WORK_DIR/$JAR "

#export LOIRA_TOKEN="<TOKEN>"

start () {
  echo "Starting loira..."
  if [ ! -f $WORK_DIR/loira.pid ]; then
    start-stop-daemon --start --verbose --background --chdir $WORK_DIR --exec $DAEMON --pidfile $WORK_DIR/loira.pid --chuid "$USER" --make-pidfile -- $DAEMON_ARGS
  else
    echo "loira is already running..."
  fi
}

stop () {
    echo "Stopping loira..."
    start-stop-daemon --stop --exec $DAEMON --pidfile $WORK_DIR/loira.pid
    rm $WORK_DIR/loira.pid
}

case $1 in
    start)
        start
    ;;
    stop)
        stop
    ;;
    restart)
        stop
        start
    ;;
esac
